# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/local/sbin/php-fpm
priority=1

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=2
